# .install/config/crontab/live
0	*	*	*	*	root	(echo "\n----------" && date && echo "----------" && bash /home/www/__dwsmyapp__.test/www/bin/console.sh dump) >> /home/www/__dwsmyapp__.test/www/var/log/dump.cron.log
