# .install/config/crontab/live
0	5	*	*	*	root	(echo "\n----------" && date && echo "----------" && bash /home/www/__dwsmyapp__.test/www/bin/console.sh mysql/dump) >> /home/www/__dwsmyapp__.test/www/var/log/mysql-dump.cron.log
