# install/config/crontab/staging
0	5	*	*	*	root	bash /home/www/myapp.test/www/bin/console.sh mysql/dump > /home/www/myapp.test/www/var/log/mysql-dump.cron.log
